# recipes-local/etckeeper/files/create-etc.service
[Unit]
Description=create etc git remote for etckeeper
Before=update-etc.service neutrino.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/sh /etc/etckeeper/create_etc.sh

[Install]
WantedBy=multi-user.target
