# recipes-local/etckeeper/files/etckeeper-autocommit.timer
[Unit]
Description=etckeeper autocommit

[Timer]
OnBootSec=1h
OnUnitInactiveSec=1d

[Install]
WantedBy=basic.target
