# recipes-neutrino/neutrino/files/neutrino.service
[Unit]
Description=Neutrino HD
After=update-kernel.service etckeeper.service

[Service]
Type=simple
ExecStart=/usr/bin/neutrino
ExecStop=killall /usr/bin/neutrino
Restart=on-abort

[Install]
WantedBy=multi-user.target
