# recipes-local/vsftpd/files/vsftpd.service
[Unit]
Description=Vsftpd ftp daemon
After=network.target

[Service]
Type=forking
ExecStart=/usr/sbin/vsftpd

[Install]
WantedBy=multi-user.target

